# Phase 2/phase-10-final-project/config_innovalabs-dashboard.service
[Unit]
Description=INNOVALABS Dashboard & Remote Control
After=network.target n8n.service ollama.service
Wants=n8n.service ollama.service

[Service]
Type=simple
User=USER_ACTUAL
Group=USER_ACTUAL
WorkingDirectory=/opt/innovalabs/dashboard
Environment="PATH=/opt/innovalabs/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="INNOVALABS_HISTORIAS_DIR=/var/opt/innovalabs/historias"
Environment="N8N_URL=http://localhost:5678"
Environment="OLLAMA_URL=http://localhost:11434"
ExecStart=/opt/innovalabs/venv/bin/python3 /opt/innovalabs/dashboard/server.py --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
